# brother-ql-label-printer/rootfs/etc/services.d/brother-ql-label-printer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here
cp /data/options.json /app/config.json
cd /app && 

## Run your program
exec python3 brother_ql_web.py
